# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid

[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
; D-06: explicit retry budget and exponential backoff. --backoff accepts a
; comma-separated list of per-attempt delays (30s, 2min, 10min). After
; --tries attempts the job lands in failed_jobs. --max-jobs recycles the
; worker periodically to avoid memory bloat from long-running processes.
command=/bin/sh -lc 'php /var/www/html/artisan queue:work ${QUEUE_CONNECTION:-database} --sleep=3 --tries=3 --backoff=30,120,600 --max-time=3600 --max-jobs=1000 --timeout=300 --no-interaction'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=3600

[program:scheduler]
command=/bin/sh -lc 'php /var/www/html/artisan schedule:work --verbose --no-interaction'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=3600
